# apps/bot/telegram-bot.service
[Unit]
Description=Telegram Bot for Task Management System
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/tipa.taska.uz/apps/bot
Environment="PATH=/var/www/tipa.taska.uz/apps/bot/venv/bin:/usr/local/bin:/usr/bin:/bin"
# ВАЖНО: переменные из .env (деплой также перезаписывает unit через apps/bot/deploy.sh)
EnvironmentFile=/var/www/tipa.taska.uz/apps/bot/.env
ExecStart=/var/www/tipa.taska.uz/apps/bot/venv/bin/python /var/www/tipa.taska.uz/apps/bot/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
